# rootfs/etc/s6-overlay/s6-rc.d/claude-skills/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Claude Skills Runner
# ==============================================================================

mkdir -p /data/runs /data/skills

# Copy built-in skills if not already present
if [ ! -f /data/skills/madplan.md ]; then
  cp /app/skills/*.md /data/skills/ 2>/dev/null || true
fi

exec node /app/server.js
